# templates/systemd/kallithea.service.erb
[Unit]
Description=Kallithea SCM Web Service
After=network.target

[Service]
User=<%= @app_user %>
WorkingDirectory=<%= @app_root %>
ExecStart=<%= @app_root %>/venv/bin/paster serve <%= @app_root %>/kallithea.ini
Restart=on-failure

[Install]
WantedBy=multi-user.target
